# systemd/pmx-lp-mini-init.service.in
[Unit]
Description=Activate Session Mode in Launchpad Mini
Before=pmx-slp-ctrl.service
After=network.target

[Service]
Type=notify
ExecStart=@BINDIR@/pmx-lp-mini-init
ExecStop=
Restart=on-failure

[Install]
WantedBy=performance-mixer.target
