# systemd/pmx-osc-sl-network-receiver.service.in
[Unit]
Description=Receives osc messages via network and send them out on a port
Before=
After=pmx-metadata-manager.service
Requires=pipewire.service wireplumber.service pmx-metadata-manager.service
Wants=

[Service]
Type=notify
ExecStart=@BINDIR@/pmx-osc-network-receiver pmx-osc-sl-network-receiver pmx-from-sl-osc 33334
ExecStop=
Restart=on-failure

[Install]
WantedBy=performance-mixer.target
Also=
