# systemd/pmx-osc-sl-network-sender.service.in
[Unit]
Description=Main OSC Network Sender for PMX-1
Before=
After=pmx-metadata-manager.service
Requires=pipewire.service wireplumber.service pmx-metadata-manager.service
Wants=

[Service]
Type=notify
ExecStart=@BINDIR@/pmx-osc-network-sender pmx-osc-sl-network-sender pmx-to-sl-osc 127.0.0.1 3300
ExecStop=
Restart=on-failure

[Install]
WantedBy=performance-mixer.target
Also=
