# systemd/pmx-sooperlooper.service.in
[Unit]
Description=Sooperlooper for PMX
Before=
After=network.target
Requires=pipewire.service wireplumber.service
Wants=

[Service]
Type=simple
ExecStart=@BINDIR@/pmx-start-sooperlooper.fish
Restart=on-failure

[Install]
WantedBy=performance-mixer.target
Also=
