# systemd/pmx-filter-chain-ctrl.service.in
[Unit]
Description=Receives osc messages and controls the filter chains used by pmx
Before=
After=pmx-metadata-manager.service
Requires=pipewire.service wireplumber.service pmx-metadata-manager.service
Wants=

[Service]
Type=notify
ExecStart=@BINDIR@/pmx-filter-chain-ctrl
ExecStop=
Restart=on-failure

[Install]
WantedBy=performance-mixer.target
Also=
