# systemd/pmx-midi-router.service.in
[Unit]
Description=Receives midi messages and sends osc for pmx
Before=
Requires=pipewire.service wireplumber.service pmx-metadata-manager.service
Wants=

[Service]
Type=notify
ExecStart=@BINDIR@/pmx-midi-router
ExecStop=
Restart=on-failure

[Install]
WantedBy=performance-mixer.target
Also=
